# {{cookiecutter.repo_name}}/config/celery.conf
; https://github.com/celery/celery/blob/3.1/extra/supervisord/celeryd.conf

[program:celery]
; Set full path to celery program if using virtualenv
environment=DJANGO_SETTINGS_MODULE="{{cookiecutter.package_name}}.settings.production",DB_USER="{{cookiecutter.package_name}}",DB_NAME="{{cookiecutter.package_name}}"
command=/var/www/{{cookiecutter.package_name}}/.venv/bin/celery worker -A {{cookiecutter.package_name}} --loglevel=INFO -E

directory=/var/www/{{cookiecutter.package_name}}
user={{cookiecutter.package_name}}
numprocs=1
stdout_logfile=/var/log/celery/worker.log
stderr_logfile=/var/log/celery/worker.log
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

; Set Celery priority higher than default (999)
; so, if rabbitmq is supervised, it will start first.
priority=1000
